# scripts/install/sh.nabla.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>sh.nabla.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/nabla-daemon</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>

  <key>ExitTimeOut</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/USERNAME/Library/Logs/nabla-daemon.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/USERNAME/Library/Logs/nabla-daemon.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NABLA_DAEMON_HOST</key>
    <string>127.0.0.1</string>
    <key>NABLA_DAEMON_PORT</key>
    <string>7777</string>
  </dict>
</dict>
</plist>
